# src/soundbox_cache.service
[Unit]
Description=Soundbox Greeting Cache Service
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/YOUR_USERNAME/soundbox
ExecStart=/home/YOUR_USERNAME/soundbox/venv/bin/python cache.py
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Timer]
OnCalendar=hourly
AccuracySec=1m
Persistent=true

[Install]
WantedBy=default.target
